# src/variations/frankenphp/etc/periodic/15min/sample
#!/bin/sh
# date | tee /tmp/test-15min.txt
